# gphile.supervisor.conf
[program:gphile]
command=/app/gphile.sh
redirect_stderr=false
